# rootfs/usr/etc/profile.d/brew.sh
#!/bin/sh
[[ -d /home/linuxbrew/ && $- == *i* ]] && eval "$(/home/linuxbrew/bin/brew shellenv)" 
