# rootfs/usr/etc/profile.d/celestite-fastfetch.sh
alias fastfetch="/usr/bin/fastfetch --config /usr/share/celestite/minimal.jsonc"
alias neofetch=fastfetch
