# rootfs/usr/etc/profile.d/rolls-fastfetch.sh
alias fastfetch="/usr/bin/fastfetch --config /usr/share/firstboot/minimal.jsonc"
alias neofetch=fastfetch
