# rootfs/usr/etc/profile.d/rolls-fastfetch.sh
alias fastfetch="/usr/bin/fastfetch --file /usr/share/firstboot/rolls.txt --logo-type file --logo-color-1 yellow"
